# resources/pkg/scripts/postinstall
#!/bin/bash
# AgentHandover post-install script
set -e

# Track non-fatal warnings that shouldn't abort install but should be reported.
WARNINGS=""
warn() { WARNINGS="${WARNINGS}\n  ⚠ $1"; }

# ----------------------------------------------------------------------
# Detect the real logged-in user (the one we want to install into).
#
# This is harder than it looks. The installer runs as root, so $USER is
# always "root" (or whatever the invoking shell was). We need "the human
# behind the keyboard" — the person whose ~/Library/LaunchAgents we
# should write to, whose ~/Library/Application Support we should create.
#
# v0.2.3 used `stat -f '%Su' /dev/console` with a fallback to `$USER`.
# That broke for at least one user (hikoae, issue #1 follow-up) because
# on his machine `/dev/console` resolved to `root` at install time,
# causing the whole install to target `/var/root/Library/LaunchAgents/`
# instead of his user's — the LaunchAgent was never visible to his
# logged-in session, `agenthandover start` failed with
# `launchctl: Load failed: 5: Input/output error`, and `doctor`
# reported `Worker launchd plist` as missing. The failure was silent
# because `/var/root` is a real directory on macOS, so the fallback
# `[ -d "$USER_HOME" ]` check always passed.
#
# v0.2.4 cascades through four methods, each more reliable than the
# previous, with explicit UID >= 500 validation and a hard-fail if none
# of them find a real user. Order matters: scutil is the authoritative
# source for "current GUI console user" per Apple's own docs; SUDO_USER
# catches the case where someone runs `sudo installer -pkg` from an
# SSH session (where scutil may return `loginwindow` or be empty);
# stat /dev/console is the legacy method and still works on most
# machines; dscl is a last-resort "first real user on this system"
# fallback for pathological cases.
# ----------------------------------------------------------------------

_uid_of() {
    # Echo the UID for a given username, or nothing if the user doesn't
    # exist or has no UID. Safe on system accounts (returns their UID,
    # which the caller validates against >= 500).
    id -u "$1" 2>/dev/null || true
}

_valid_user() {
    # Return 0 if $1 is a non-empty, non-"root", non-"loginwindow",
    # non-"_appstore" (system account) username with a UID >= 500.
    local _u="$1"
    if [ -z "$_u" ] || [ "$_u" = "root" ] || [ "$_u" = "loginwindow" ]; then
        return 1
    fi
    local _uid=$(_uid_of "$_u")
    if [ -z "$_uid" ]; then
        return 1
    fi
    if [ "$_uid" -lt 500 ] 2>/dev/null; then
        return 1
    fi
    return 0
}

CONSOLE_USER=""

# Method 1: scutil on State:/Users/ConsoleUser — Apple's recommended
# source for the current GUI user. Returns empty or "loginwindow" if
# no one is logged in at the console.
_method1=$(scutil <<< "show State:/Users/ConsoleUser" 2>/dev/null | awk '/Name :/ && $3 != "loginwindow" { print $3; exit }')
if _valid_user "${_method1}"; then
    CONSOLE_USER="${_method1}"
    _detect_method="scutil"
fi

# Method 2: SUDO_USER env var — reliable when installed via
# `sudo installer -pkg ...` from a terminal, even if no one is at the
# GUI console. Installer.app invocations don't set this (they run as
# root directly), but `sudo installer` CLI invocations do.
if [ -z "${CONSOLE_USER}" ]; then
    if _valid_user "${SUDO_USER}"; then
        CONSOLE_USER="${SUDO_USER}"
        _detect_method="SUDO_USER"
    fi
fi

# Method 3: stat -f '%Su' /dev/console — the legacy v0.2.3 method.
# Kept as a fallback because it works on most machines, even though
# it broke for hikoae.
if [ -z "${CONSOLE_USER}" ]; then
    _method3=$(stat -f '%Su' /dev/console 2>/dev/null)
    if _valid_user "${_method3}"; then
        CONSOLE_USER="${_method3}"
        _detect_method="stat-console"
    fi
fi

# Method 4: first real user on the system via dscl. Last resort —
# picks an arbitrary user on multi-user machines but is better than
# targeting /var/root.
if [ -z "${CONSOLE_USER}" ]; then
    _method4=$(dscl . list /Users UniqueID 2>/dev/null | awk '$2 >= 500 && $2 < 60000 { print $1; exit }')
    if _valid_user "${_method4}"; then
        CONSOLE_USER="${_method4}"
        _detect_method="dscl-first-real-user"
    fi
fi

# Hard fail if every method failed. Better than silently installing
# into /var/root and leaving the user with a broken install.
if [ -z "${CONSOLE_USER}" ]; then
    cat >&2 <<'ERR'
ERROR: AgentHandover postinstall could not determine which user to install for.

All four detection methods returned empty or invalid:
  1. scutil 'State:/Users/ConsoleUser'  (the GUI console user)
  2. SUDO_USER env var                    (the user who ran sudo)
  3. stat -f '%Su' /dev/console           (the console device owner)
  4. dscl first-real-user scan            (any user with UID >= 500)

This usually happens when:
  - No one is logged in at the GUI console (SSH-only install without sudo)
  - The install ran before the user's login session fully initialized

To fix: log in to the Mac via the GUI, then re-run:
    sudo installer -pkg AgentHandover-*.pkg -target /

If you're installing via SSH from a remote machine, make sure you're
using `sudo` (not `sudo -i` or `su -`) so SUDO_USER is set correctly.
ERR
    exit 1
fi

CONSOLE_UID=$(_uid_of "${CONSOLE_USER}")
USER_HOME=$(dscl . -read /Users/"${CONSOLE_USER}" NFSHomeDirectory 2>/dev/null | awk '{print $2}')
if [ -z "$USER_HOME" ] || [ ! -d "$USER_HOME" ]; then
    USER_HOME=$(eval echo ~"$CONSOLE_USER")
fi
# Safety: if USER_HOME somehow resolved to /var/root (i.e. user is
# root despite the _valid_user check — impossible, but let's be sure)
# or is still empty, hard fail.
if [ -z "$USER_HOME" ] || [ ! -d "$USER_HOME" ] || [ "$USER_HOME" = "/var/root" ]; then
    echo "ERROR: AgentHandover postinstall resolved CONSOLE_USER='${CONSOLE_USER}' to USER_HOME='${USER_HOME}' which is invalid." >&2
    echo "Refusing to install LaunchAgent into system user home directory." >&2
    exit 1
fi

echo "  Installing for user: ${CONSOLE_USER} (UID=${CONSOLE_UID}, HOME=${USER_HOME}, detected via ${_detect_method})"
DATA_DIR="${USER_HOME}/Library/Application Support/agenthandover"
LAUNCH_AGENTS="${USER_HOME}/Library/LaunchAgents"
LIB_DIR="/usr/local/lib/agenthandover"

# Create data directories (owned by the user, not root)
mkdir -p "${DATA_DIR}/logs"
mkdir -p "${DATA_DIR}/artifacts"
chown -R "${CONSOLE_USER}" "${DATA_DIR}"

# Write launchd plists with user-specific log paths
# (cannot use static plists because log paths need the real user's home)
mkdir -p "${LAUNCH_AGENTS}"
LOG_DIR="${DATA_DIR}/logs"

# Daemon is launched as an app bundle by the menu bar app (not launchd).
# Only the worker uses launchd.
cp "${LIB_DIR}/launchd/com.agenthandover.worker.plist" "${LAUNCH_AGENTS}/"
sed -i '' \
    -e "s|/tmp/com.agenthandover.worker.stdout.log|${LOG_DIR}/worker.stdout.log|" \
    -e "s|/tmp/com.agenthandover.worker.stderr.log|${LOG_DIR}/worker.stderr.log|" \
    "${LAUNCH_AGENTS}/com.agenthandover.worker.plist"

chown "${CONSOLE_USER}" "${LAUNCH_AGENTS}/com.agenthandover.worker.plist"

# Create default config if none exists
if [ ! -f "${DATA_DIR}/config.toml" ]; then
    cat > "${DATA_DIR}/config.toml" << 'TOML'
# AgentHandover Configuration
# See https://github.com/sandroandric/AgentHandover for documentation

[observer]
t_dwell_seconds = 3
t_scroll_read_seconds = 8
capture_screenshots = true
screenshot_max_per_minute = 20

[privacy]
enable_inline_secret_redaction = true
secure_field_drop = true

[storage]
retention_days_raw = 14
retention_days_episodes = 90

[export]
adapter = "openclaw"
json_export = false

[embedding]
model = "nomic-embed-text"
image_embeddings = true
TOML
    chown "${CONSOLE_USER}" "${DATA_DIR}/config.toml"
fi

# Write the Chrome native messaging host manifest for every installed
# Chromium-family browser. Previous versions only DELETED stale manifests
# here and relied on the app to recreate them on launch via
# ServiceController.installNativeMessagingHostManifest(). That broke when
# the app didn't launch cleanly after install (crash, timing, permission
# issue) — the manifests stayed deleted and the user had to run
# `agenthandover setup --extension` manually. v0.2.6 writes them directly
# in the postinstall as well (belt and suspenders — the app will overwrite
# them on launch anyway, but this guarantees they exist).
#
# The extension ID is derived from the SHA-256 of the RSA key in
# extension/manifest.json, mapped to Chrome's a-p alphabet.
NM_EXTENSION_ID="jpemkdcihaijkolbkankcldmiimmmnfo"
NM_MANIFEST='{
  "name": "com.agenthandover.host",
  "description": "AgentHandover native messaging host",
  "path": "/usr/local/lib/agenthandover/ah-observer",
  "type": "stdio",
  "allowed_origins": [
    "chrome-extension://'"${NM_EXTENSION_ID}"'/"
  ]
}'
for browser_dir in \
    "Google/Chrome" \
    "Chromium" \
    "BraveSoftware/Brave-Browser" \
    "Microsoft Edge" \
    "Arc/User Data" \
    "Comet"; do
    BROWSER_ROOT="${USER_HOME}/Library/Application Support/${browser_dir}"
    if [ -d "${BROWSER_ROOT}" ]; then
        mkdir -p "${BROWSER_ROOT}/NativeMessagingHosts"
        echo "${NM_MANIFEST}" > "${BROWSER_ROOT}/NativeMessagingHosts/com.agenthandover.host.json"
        chown "${CONSOLE_USER}" "${BROWSER_ROOT}/NativeMessagingHosts/com.agenthandover.host.json"
    fi
done

# Set up Python virtual environment for worker.
# The worker is required for VLM annotation and SOP generation — without it,
# AgentHandover cannot function. Treat failure as fatal.
if [ ! -d "${LIB_DIR}/venv" ]; then
    # Find Python 3.11+ — installer runs as root so PATH is minimal.
    # Search common locations before falling back to system python3.
    PYTHON3=""
    for candidate in \
        /opt/homebrew/bin/python3 \
        /usr/local/bin/python3 \
        /opt/homebrew/bin/python3.14 \
        /opt/homebrew/bin/python3.13 \
        /opt/homebrew/bin/python3.12 \
        /opt/homebrew/bin/python3.11 \
        /usr/local/bin/python3.13 \
        /usr/local/bin/python3.12 \
        /usr/local/bin/python3.11 \
        /Library/Frameworks/Python.framework/Versions/3.14/bin/python3 \
        /Library/Frameworks/Python.framework/Versions/3.13/bin/python3 \
        /Library/Frameworks/Python.framework/Versions/3.12/bin/python3 \
        /Library/Frameworks/Python.framework/Versions/3.11/bin/python3 \
        "$(which python3 2>/dev/null)" \
        python3; do
        if [ -x "$candidate" ] 2>/dev/null; then
            ver=$("$candidate" -c 'import sys; print(f"{sys.version_info.major}.{sys.version_info.minor}")' 2>/dev/null || echo "0.0")
            major=$(echo "$ver" | cut -d. -f1)
            minor=$(echo "$ver" | cut -d. -f2)
            if [ "$major" -ge 3 ] 2>/dev/null && [ "$minor" -ge 11 ] 2>/dev/null; then
                PYTHON3="$candidate"
                break
            fi
        fi
    done

    if [ -z "$PYTHON3" ]; then
        echo "ERROR: Python 3.11+ not found. AgentHandover requires Python 3.11 or later."
        echo "Install Python from https://www.python.org/downloads/ or via: brew install python@3.13"
        exit 1
    fi

    echo "  Using Python: $PYTHON3"
    PYTHON_VERSION=$("$PYTHON3" -c 'import sys; print(f"{sys.version_info.major}.{sys.version_info.minor}")' 2>/dev/null)
    echo "  Python version: ${PYTHON_VERSION}"

    # Create venv as the real user (not root) so future pip installs
    # don't require sudo. The installer runs as root but the venv
    # should be owned by the console user.
    chown -R "${CONSOLE_USER}" "${LIB_DIR}"
    if ! sudo -u "${CONSOLE_USER}" "$PYTHON3" -m venv "${LIB_DIR}/venv"; then
        echo "ERROR: Failed to create Python virtual environment at ${LIB_DIR}/venv."
        echo "Try manually: $PYTHON3 -m venv ${LIB_DIR}/venv"
        exit 1
    fi

    if ! sudo -u "${CONSOLE_USER}" "${LIB_DIR}/venv/bin/pip" install --quiet "${LIB_DIR}/worker/"; then
        echo "ERROR: Failed to install AgentHandover worker Python package."
        echo "Try manually: ${LIB_DIR}/venv/bin/pip install ${LIB_DIR}/worker/"
        rm -rf "${LIB_DIR}/venv"
        exit 1
    fi

    # Install optional dependencies for image embeddings (Apple Silicon only)
    sudo -u "${CONSOLE_USER}" "${LIB_DIR}/venv/bin/pip" install --quiet "mlx-embeddings>=0.0.4" 2>/dev/null || true

    # Ensure the venv is owned by the user, not root
    chown -R "${CONSOLE_USER}" "${LIB_DIR}/venv"

    echo "  Python virtual environment created."
else
    # Venv exists — reinstall worker package to pick up code changes.
    # chown first: earlier installs may have left root-owned dist-info
    # which blocks pip uninstall/reinstall as the console user.
    chown -R "${CONSOLE_USER}" "${LIB_DIR}/venv"
    echo "  Updating worker package..."
    sudo -u "${CONSOLE_USER}" "${LIB_DIR}/venv/bin/pip" install --quiet --force-reinstall "${LIB_DIR}/worker/" 2>/dev/null || true
    # Ensure mlx-embeddings is installed
    sudo -u "${CONSOLE_USER}" "${LIB_DIR}/venv/bin/pip" install --quiet "mlx-embeddings>=0.0.4" 2>/dev/null || true
fi

# Set proper permissions
chmod 755 /usr/local/lib/agenthandover/ah-observer
chmod 755 /usr/local/bin/agenthandover

# Symlink worker entry points to /usr/local/bin/ so agents and CLI can find them.
# These are Python scripts installed in the venv — the symlinks make them PATH-accessible.
for cmd in agenthandover-mcp agenthandover-connect agenthandover-worker; do
    if [ -x "${LIB_DIR}/venv/bin/${cmd}" ]; then
        ln -sf "${LIB_DIR}/venv/bin/${cmd}" "/usr/local/bin/${cmd}"
    fi
done

# Clear stale daemon TCC identity so users only ever see AgentHandover in
# System Settings going forward. This does not reset the main app grant.
sudo -u "${CONSOLE_USER}" tccutil reset ScreenCapture com.agenthandover.daemon 2>/dev/null || true
sudo -u "${CONSOLE_USER}" tccutil reset Accessibility com.agenthandover.daemon 2>/dev/null || true
sudo -u "${CONSOLE_USER}" launchctl disable "gui/${CONSOLE_UID}/com.agenthandover.daemon" 2>/dev/null || true

# Only reset onboarding on FRESH install (data dir didn't exist before)
if [ ! -f "${DATA_DIR}/config.toml.bak" ]; then
    sudo -u "${CONSOLE_USER}" defaults write com.agenthandover.app hasCompletedOnboarding -bool false 2>/dev/null || true
fi
# Mark that install has run (for upgrade detection)
touch "${DATA_DIR}/config.toml.bak" 2>/dev/null || true

# Chrome extension setup (informational only — never fatal)
EXT_DIR="${LIB_DIR}/extension"
if [ -f "${EXT_DIR}/manifest.json" ] && [ -f "${EXT_DIR}/background.js" ]; then
    echo "Chrome extension pre-built and ready."
    echo "  Load it in Chrome: chrome://extensions -> Load unpacked -> ${EXT_DIR}"
elif [ -d "${EXT_DIR}/src" ]; then
    echo ""
    echo "Chrome extension source is included but not pre-built."
    echo "To build the extension:"
    echo "  cd ${EXT_DIR}"
    echo "  npm install && npm run build"
    echo "Then load in Chrome: chrome://extensions -> Load unpacked -> ${EXT_DIR}/dist"
else
    echo ""
    echo "Chrome extension not included in this package."
    echo "To install it from source:"
    echo "  git clone https://github.com/sandroandric/AgentHandover.git"
    echo "  cd AgentHandover/extension && npm install && npm run build"
    echo "Then load in Chrome: chrome://extensions -> Load unpacked -> extension/dist"
fi

# Register launchd services (RunAtLoad=false, so they won't start yet)
# Only the worker uses launchd now. The daemon is started directly by the app.
# Use modern launchctl bootstrap (load/unload deprecated on macOS 13+)
sudo -u "${CONSOLE_USER}" launchctl bootstrap "gui/${CONSOLE_UID}" "${LAUNCH_AGENTS}/com.agenthandover.worker.plist" 2>/dev/null || true
echo "  Services registered (will start when app launches)."

# Launch the menu bar app (as the real user, not root)
if [ -d "/Applications/AgentHandover.app" ]; then
    sudo -u "${CONSOLE_USER}" open /Applications/AgentHandover.app 2>/dev/null || true
    echo "AgentHandover app launched — check your menu bar."
fi

echo ""
if [ -n "$WARNINGS" ]; then
    echo "AgentHandover installed with warnings:"
    echo -e "$WARNINGS"
    echo ""
    echo "Run 'agenthandover doctor' to diagnose issues."
    exit 0
else
    echo "AgentHandover installed successfully."
    exit 0
fi
